# infra/aws/systemd/kalshi-observability.timer
[Unit]
Description=Run Kalshi Tier 2 observability after the daily health report

[Timer]
OnCalendar=*-*-* 06:25:00
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
