# infra/aws/systemd/kalshi-forward.service
[Unit]
Description=Kalshi forward ingestion (update_forward.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/kalshi-pipeline
# Optional: UPDATE_FORWARD_EXTRA_ARGS (e.g. --markets-only). See forward-exec.env.example.
EnvironmentFile=-/etc/kalshi/forward-exec.env
# Wrapper avoids systemd interpreting ${...} in ExecStart; logs show script path in journal.
ExecStart=/opt/kalshi-pipeline/infra/aws/run-update-forward.sh
TimeoutStartSec=0
Nice=5

[Install]
WantedBy=multi-user.target
