# infra/aws/systemd/kalshi-health.service
[Unit]
Description=Kalshi data health validation (validate_data_health.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/kalshi-pipeline
ExecStartPre=/bin/install -d -m 755 -o ubuntu -g ubuntu /opt/kalshi-pipeline/data/kalshi/state
ExecStart=/bin/bash -lc 'set -a && [ -f /opt/kalshi-pipeline/.env ] && . /opt/kalshi-pipeline/.env && set +a && exec /usr/local/bin/uv run python scripts/validate_data_health.py --output /opt/kalshi-pipeline/data/kalshi/state/health_report.json'
TimeoutStartSec=0
Nice=10

[Install]
WantedBy=multi-user.target
