# infra/aws/systemd/kalshi-forward.service
[Unit]
Description=Kalshi forward ingestion (update_forward.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/kalshi-pipeline
# Extra args come from /etc/kalshi/forward-exec.env inside run-update-forward.sh (not EnvironmentFile= here).
ExecStart=/opt/kalshi-pipeline/infra/aws/run-update-forward.sh
TimeoutStartSec=0
Nice=5

[Install]
WantedBy=multi-user.target
